# .config/systemd/user/elephant.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY
Description=elephant

[Service]
Type=exec
ExitType=cgroup
Slice=app-graphical.slice
ExecStart=%h/.local/bin/elephant
Restart=on-failure
TimeoutStopFailureMode=terminate

[Install]
WantedBy=graphical-session.target
