# .config_wm/systemd/user/container-ollama.service
# container-ollama.service
# autogenerated by Podman 5.7.1
# Wed Jan 14 01:31:37 CST 2026

[Unit]
Description=Podman container-ollama.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/user/1000/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start ollama
ExecStop=/usr/bin/podman stop  \
	-t 10 ollama
ExecStopPost=/usr/bin/podman stop  \
	-t 10 ollama
PIDFile=/run/user/1000/containers/overlay-containers/5b0a11b4dfc31c6146b1c147ee5c0839e1e74ee2bcf9535d097e7d6e75c19ef4/userdata/conmon.pid
Type=forking

[Install]
WantedBy=default.target
